# .crontab
0 * * * * php -f /path/to/splitter/feed.php > /dev/null